# src/terminal-session-snapshot.timer
[Unit]
Description=Auto-snapshot terminal session every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
